# deploy/systemd/yoyopod-slot.service
[Unit]
Description=YoyoPod (slot deploy, OTA-ready)
After=network-online.target sound.target pisugar-server.service
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3
OnFailure=yoyopod-rollback.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/yoyopod-slot
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/
ExecStart=/usr/bin/env bash -lc 'exec "$${YOYOPOD_ROOT:-/opt/yoyopod}/current/bin/launch"'
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
