# deploy/systemd/yoyopod-rollback.service
[Unit]
Description=YoyoPod rollback after repeated slot failures
StartLimitIntervalSec=3600
StartLimitBurst=2

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/yoyopod-slot
ExecStart=/usr/bin/env bash -lc 'exec "$${YOYOPOD_ROOT:-/opt/yoyopod}/bin/rollback.sh"'
StandardOutput=journal
StandardError=journal
