# deploy/systemd/yoyopod-dev.service
[Unit]
Description=YoyoPod dev lane (mutable checkout)
After=network-online.target sound.target pisugar-server.service
Wants=network-online.target
Conflicts=yoyopod-prod.service

[Service]
Type=simple
EnvironmentFile=-/etc/default/yoyopod-dev
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/root
WorkingDirectory=/
ExecStartPre=/usr/bin/env bash -lc 'HOME="$${HOME:-/root}"; install -d -m 700 "$$HOME/.local/share/linphone" "$$HOME/.config/linphone"; CHECKOUT="$${YOYOPOD_DEV_CHECKOUT:-/opt/yoyopod-dev/checkout}"; test -f "$$CHECKOUT/pyproject.toml"; test -x "$$CHECKOUT/device/runtime/build/yoyopod-runtime"; test -x "$$CHECKOUT/device/ui/build/yoyopod-ui-host"'
ExecStart=/usr/bin/env bash -lc 'CHECKOUT="$${YOYOPOD_DEV_CHECKOUT:-/opt/yoyopod-dev/checkout}"; cd "$$CHECKOUT"; exec "$$CHECKOUT/device/runtime/build/yoyopod-runtime" --config-dir "$$CHECKOUT/config" --hardware whisplay'
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
