# deploy/systemd/yoyopod@.service
[Unit]
Description=YoYoPod portable VoIP/music appliance (%i)
After=network-online.target sound.target pisugar-server.service
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/etc/default/yoyopod
ExecStartPre=/usr/bin/env bash -lc 'PROJECT_DIR="${YOYOPOD_PROJECT_DIR:-$HOME/yoyopod-core}"; case "$PROJECT_DIR" in "~") PROJECT_DIR="$HOME" ;; "~/"*) PROJECT_DIR="$HOME/${PROJECT_DIR#~/}" ;; esac; if [ ! -f "$PROJECT_DIR/pyproject.toml" ] && [ -f "$HOME/YoyoPod_Core/pyproject.toml" ]; then PROJECT_DIR="$HOME/YoyoPod_Core"; fi; if [ ! -f "$PROJECT_DIR/pyproject.toml" ] && [ -f "$HOME/yoyo-py/pyproject.toml" ]; then PROJECT_DIR="$HOME/yoyo-py"; fi; if [ ! -f "$PROJECT_DIR/pyproject.toml" ]; then echo "Could not find a YoYoPod checkout at $PROJECT_DIR" >&2; exit 1; fi; PYTHON_BIN="$PROJECT_DIR/.venv/bin/python"; if [ ! -x "$PYTHON_BIN" ]; then echo "Could not find a YoYoPod virtualenv python at $PYTHON_BIN" >&2; exit 1; fi; cd "$PROJECT_DIR"; exec "$PYTHON_BIN" -m yoyopod_cli.main build ensure-native'
ExecStart=/usr/bin/env bash -lc 'PROJECT_DIR="${YOYOPOD_PROJECT_DIR:-$HOME/yoyopod-core}"; case "$PROJECT_DIR" in "~") PROJECT_DIR="$HOME" ;; "~/"*) PROJECT_DIR="$HOME/${PROJECT_DIR#~/}" ;; esac; if [ ! -f "$PROJECT_DIR/pyproject.toml" ] && [ -f "$HOME/YoyoPod_Core/pyproject.toml" ]; then PROJECT_DIR="$HOME/YoyoPod_Core"; fi; if [ ! -f "$PROJECT_DIR/pyproject.toml" ] && [ -f "$HOME/yoyo-py/pyproject.toml" ]; then PROJECT_DIR="$HOME/yoyo-py"; fi; if [ ! -f "$PROJECT_DIR/pyproject.toml" ]; then echo "Could not find a YoYoPod checkout at $PROJECT_DIR" >&2; exit 1; fi; PYTHON_BIN="$PROJECT_DIR/.venv/bin/python"; if [ ! -x "$PYTHON_BIN" ]; then echo "Could not find a YoYoPod virtualenv python at $PYTHON_BIN" >&2; exit 1; fi; cd "$PROJECT_DIR"; exec "$PYTHON_BIN" yoyopod.py'
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
